# systemd/system/foo-watchdog.service
[Unit]
Description=Foo Watchdog

[Service]
ExecStart=/root/watchdog

[Install]
WantedBy=multi-user.target

